# scripts/templates/registry.service
[Unit]
Description=Subdomain Registry Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www
ExecStart=/usr/local/bin/bun run /var/www/registry-server.ts
Restart=always
RestartSec=10

# Environment variables (set these in /etc/registry.env)
EnvironmentFile=/etc/registry.env

[Install]
WantedBy=multi-user.target
